# config/systemd/perfect10k-docker.service
[Unit]
Description=Perfect10k API Server (Docker Container)
After=docker.service network.target
Requires=docker.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
WorkingDirectory=/var/www/Perfect10k

# Ensure Docker is accessible
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do sleep 1; done'

# Stop any existing container first
ExecStartPre=-/usr/bin/docker compose down

# Build only if image doesn't exist
ExecStartPre=/bin/bash -c 'docker image inspect perfect10k-perfect10k-app >/dev/null 2>&1 || docker compose build'

# Start the container
ExecStart=/usr/bin/docker compose up -d

# Stop the container
ExecStop=/usr/bin/docker compose down

# Health check
ExecStartPost=/bin/bash -c 'sleep 10 && curl -f http://localhost:8000/health || exit 1'

# Restart policy
Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perfect10k-docker

[Install]
WantedBy=multi-user.target