# config/systemd/perfect10k-api.service
[Unit]
Description=Perfect10k API Server
After=network.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/Perfect10k/backend

# Set up the environment with the virtual environment activated
Environment=PATH=/var/www/Perfect10k/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/var/www/Perfect10k/.venv
Environment=PYTHONPATH=/var/www/Perfect10k
Environment=LOG_DIR=/var/www/Perfect10k/logs

# Use the Python from the virtual environment directly instead of uv
ExecStart=/var/www/Perfect10k/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1

# Restart policy
Restart=on-failure
RestartSec=30
TimeoutStopSec=30
KillMode=mixed

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/Perfect10k/backend/cache
ReadWritePaths=/var/www/Perfect10k/cache
ReadWritePaths=/var/www/Perfect10k/logs
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perfect10k-api

[Install]
WantedBy=multi-user.target