# config/systemd/perfect10k-docker.service
[Unit]
Description=Perfect10k API Server (Docker Container)
After=docker.service network.target
Requires=docker.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
WorkingDirectory=/var/www/Perfect10k

# Ensure Docker is accessible
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do sleep 1; done'

# Stop and remove any existing container
ExecStartPre=-/usr/bin/docker stop perfect10k-app
ExecStartPre=-/usr/bin/docker rm perfect10k-app

# Create necessary directories
ExecStartPre=/bin/mkdir -p /var/www/Perfect10k/logs
ExecStartPre=/bin/mkdir -p /var/www/Perfect10k/backend/cache
ExecStartPre=/bin/mkdir -p /var/www/Perfect10k/backend/storage

# Start the container
ExecStart=/usr/bin/docker run -d \
    --name perfect10k-app \
    --restart unless-stopped \
    -p 8000:8000 \
    -v /var/www/Perfect10k/logs:/app/logs \
    -v /var/www/Perfect10k/backend/cache:/app/backend/cache \
    -v /var/www/Perfect10k/backend/storage:/app/backend/storage \
    perfect10k-app:latest

# Stop the container
ExecStop=/usr/bin/docker stop perfect10k-app
ExecStop=/usr/bin/docker rm perfect10k-app

# Health check removed - service takes too long to start for fixed timeout

# Restart policy
Restart=on-failure
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perfect10k-docker

[Install]
WantedBy=multi-user.target